# deploy/systemd/ivyea-ops.service.template
[Unit]
Description=IvyeaOps (Amazon operations hub)
After=network.target

[Service]
Type=simple
User=${IVYEA_OPS_USER}
WorkingDirectory=${INSTALL_DIR}/server
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=${INSTALL_DIR}/server
ExecStart=${PYTHON_BIN} -m uvicorn app.main:app --host 127.0.0.1 --port ${IVYEA_OPS_PORT} --log-level info
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
